# deployment/vps/shipit-updater.path
[Unit]
Description=Watch for ShipIt update requests

[Path]
PathExists=/opt/shipit/.update-requested
Unit=shipit-updater.service

[Install]
WantedBy=multi-user.target
