# deployment/hetzner/shipit-updater.service
[Unit]
Description=ShipIt self-updater
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/shipit/deployment/hetzner/update.sh
StandardOutput=journal
StandardError=journal
